# skills/agent-chat/templates/systemd/everclaw-agent-chat.service
[Unit]
Description=EverClaw XMTP Agent Chat Daemon
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory={{EVERCLAW_PATH}}
ExecStart={{NODE_BIN}} {{DAEMON_PATH}}
Restart=on-failure
RestartSec=10
Environment=NODE_ENV=production

# Logs to user journal (view with: journalctl --user -u everclaw-agent-chat)
StandardOutput=journal
StandardError=journal
SyslogIdentifier=everclaw-agent-chat

[Install]
WantedBy=default.target